# mothics.service
[Unit]
Description=Mothics CLI Service
After=network.target

[Service]
User=root
WorkingDirectory=/root/central-unit
ExecStart=/bin/bash /root/central-unit/cli.sh
ExecStop=/usr/bin/tmux kill-session -t mothics
Restart=on-failure
Type=forking

[Install]
WantedBy=multi-user.target

