# mothics.service
[Unit]
Description=Mothics CLI Service
After=network.target

[Service]
User=%i
WorkingDirectory=/home/%i/mothics
ExecStart=/bin/bash /home/%i/mothics/cli.sh
ExecStop=/usr/bin/tmux kill-session -t mothics
Restart=on-failure
Type=forking

[Install]
WantedBy=multi-user.target

